# collector/com.usage-badge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.usage-badge.collector</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_PYTHON</string>
    <string>REPLACE_HOME/usage-badge/collector/collect.py</string>
  </array>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/tmp/usage-badge.log</string>
  <key>StandardErrorPath</key><string>/tmp/usage-badge.log</string>
</dict>
</plist>
