# collector/com.usage-badge.ollama-meter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.usage-badge.ollama-meter</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_PYTHON</string>
    <string>REPLACE_HOME/usage-badge/collector/ollama_meter.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Point at your real Ollama. Local daemon by default; use
         https://ollama.com for Ollama Cloud. -->
    <key>OLLAMA_TARGET</key><string>http://localhost:11434</string>
    <key>METER_PORT</key><string>11435</string>
  </dict>
  <key>KeepAlive</key><true/>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/tmp/usage-badge-ollama-meter.log</string>
  <key>StandardErrorPath</key><string>/tmp/usage-badge-ollama-meter.log</string>
</dict>
</plist>
